# examples/systemd/ExampleService.conf
start on net-device-up
exec /home/webmaster/new-site/www/app